# vendor/ruby/1.8/gems/chef-0.9.12/distro/debian/etc/init/chef-solr.conf
# chef-solr - Chef Solr
#
# Chef Solr provides Solr search for Chef Server

description "Chef Solr"

start on filesystem
stop on runlevel [!2345]

respawn
respawn limit 5 30

pre-start script
    test -x /usr/bin/chef-solr || { stop; exit 0; }
end script

exec /usr/bin/chef-solr -c /etc/chef/solr.rb -L /var/log/chef/solr.log
