# vendor/ruby/1.8/gems/chef-0.9.12/distro/debian/etc/init/chef-server.conf
# chef-server - Chef Server
#
# Chef Server provides the Chef API server

description "Chef Server API"

start on filesystem
stop on runlevel [!2345]

respawn
respawn limit 5 30

pre-start script
    test -x /usr/bin/chef-server || { stop; exit 0; }
end script

exec /usr/bin/chef-server -e production -p 4000 -L /var/log/chef/server.log
